# nodejs.conf
[program:nodejs]
command=nodejs /usr/local/nodejs/main.js
autostart=true
autorestart=true

